# ops/systemd/fangabase-mail-worker.service
[Unit]
Description=FangaBase transactional mail worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/fangabase/apps/server
ExecStart=/usr/bin/php artisan fangabase:mail-worker --limit=25
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
User=www-data
Group=www-data
EnvironmentFile=/etc/fangabase/server.env
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
